# button-capture.service
[Unit]
Description=button capture service
After=home.mount

[Service]
ExecStart=/opt/bin/button-capture

[Install]
WantedBy=multi-user.target
